# data/etc/periodic/15min/mirror
#!/bin/sh

git svn fetch
git push origin
