# data/etc/periodic/15min/mirror
#!/bin/sh

cd /work
git svn fetch
git push origin
